# network-monitor.service
[Unit]
Description=Network Monitoring System
After=network.target

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/network-monitor
Environment="PATH=/opt/network-monitor/venv/bin"
ExecStart=/opt/network-monitor/venv/bin/python /opt/network-monitor/api.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=network-monitor

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
